# com.kairos.menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.kairos.menubar.plist — macOS launchd agent (auto-start on login)
  =====================================================================
  Optional: use this to auto-start Kairos Agent when you log in,
  instead of adding it to Login Items manually.

  Setup:
    1. Edit the two paths below to match where you cloned kairos_agents/
    2. Copy to ~/Library/LaunchAgents/com.kairos.menubar.plist
    3. Run: launchctl load ~/Library/LaunchAgents/com.kairos.menubar.plist

  To stop auto-start:
    launchctl unload ~/Library/LaunchAgents/com.kairos.menubar.plist
    rm ~/Library/LaunchAgents/com.kairos.menubar.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kairos.menubar</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Edit this path to where you cloned kairos_agents/ -->
        <string>/path/to/kairos_agents/run.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <!-- Edit this path to match above -->
    <string>/path/to/kairos_agents</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/kairos-launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/kairos-launchd.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
